# yawgi/meta/yawgi.service
[Service]
ExecStart=/home/watduhhekbro/external/tmp/yawgi/yawgi

[Install]
WantedBy=multi-user.target
